# roles/unbound/files/unbound.upstart
description "unbound - Validating, recursive, and caching DNS resolver"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

exec /usr/sbin/unbound -c /etc/unbound/unbound.conf
